# launchd/com.claude-maintenance.maintenance.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.claude-maintenance.maintenance</string>
  <key>ProgramArguments</key>
  <array>
    <string>$HOME/.claude/cron-jobs/scripts/run-maintenance.sh</string>
  </array>
  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Weekday</key><integer>0</integer>
      <key>Hour</key><integer>22</integer>
      <key>Minute</key><integer>0</integer>
    </dict>
    <dict>
      <key>Weekday</key><integer>3</integer>
      <key>Hour</key><integer>22</integer>
      <key>Minute</key><integer>0</integer>
    </dict>
  </array>
  <key>RunAtLoad</key>
  <false/>
  <key>StandardOutPath</key>
  <string>$HOME/.claude/cron-jobs/reports/launchd-maintenance.out.log</string>
  <key>StandardErrorPath</key>
  <string>$HOME/.claude/cron-jobs/reports/launchd-maintenance.err.log</string>
</dict>
</plist>
